# packaging/debian/postinst
#!/bin/bash
set -e

# Reload systemd user units for the installing user (not root)
if [ -n "$SUDO_USER" ]; then
    REAL_USER="$SUDO_USER"
elif [ -n "$PKEXEC_UID" ]; then
    REAL_USER=$(id -nu "$PKEXEC_UID")
else
    REAL_USER=""
fi

if [ -n "$REAL_USER" ] && [ "$REAL_USER" != "root" ]; then
    REAL_UID=$(id -u "$REAL_USER")
    sudo -u "$REAL_USER" XDG_RUNTIME_DIR="/run/user/$REAL_UID" DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/$REAL_UID/bus" \
        systemctl --user daemon-reload || true
    sudo -u "$REAL_USER" XDG_RUNTIME_DIR="/run/user/$REAL_UID" DBUS_SESSION_BUS_ADDRESS="unix:path=/run/user/$REAL_UID/bus" \
        systemctl --user enable --now pomodoro.service || true
else
    echo "NOTE: Run these commands to start the daemon:"
    echo "  systemctl --user daemon-reload"
    echo "  systemctl --user enable --now pomodoro"
fi
